# base/packages/graphical/alsa-udev/alsa-udev.rules
SUBSYSTEM=="sound", KERNEL=="controlC*", RUN+="/lib/udev/alsa-udev %n"
